# slock/systemd/user/slock.service
[Unit]
Description=Slock Daemon
After=network-online.target slock-bridge.service
Wants=network-online.target
{{#if slock_proxy}}
Requires=slock-bridge.service
{{/if}}

[Service]
Type=simple
WorkingDirectory=%h/.config/slock
ExecStartPre=bun update --latest --force @slock-ai/daemon
{{#if slock_proxy}}
ExecStart=bun run slock-daemon --server-url http://127.0.0.1:19280 --api-key {{slock_api_key}}
{{else}}
ExecStart=bun run slock-daemon --server-url https://api.slock.ai --api-key {{slock_api_key}}
{{/if}}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
