# bub/systemd/user/bub.service
[Unit]
Description=Bub, a collaborative agent
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.config/bub
ExecStart=uv run --isolated --no-project --with git+https://github.com/bubbuild/bub.git --with git+https://github.com/bubbuild/bub-contrib.git#subdirectory=packages/bub-feishu bub gateway --enable-channel feishu
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
