# slock-ai/systemd/user/slock-ai-bridge.service
[Unit]
Description=Slock AI WebSocket Bridge (proxy tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.config/slock-ai
Environment="HTTPS_PROXY={{slock_proxy}}"
Environment="NODE_PATH=%h/.config/slock-ai/node_modules"
ExecStart=node ws-bridge.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
