# droid/systemd/droid.service
[Unit]
Description=Droid Agent Daemon with Remote Access
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/droid daemon --remote-access
{{#if agent_tunnel.wss_proxy}}
Environment=WSS_PROXY={{agent_tunnel.wss_proxy}}
{{/if}}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
