# raft/systemd/user/raft.service
[Unit]
Description=Raft.build - Where humans and AI agents build together.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.config/raft
ExecStartPre=bun update --latest --force @botiverse/raft-daemon@latest
ExecStart=bun run raft-daemon --server-url https://api.raft.build --api-key {{raft.api_key}}
{{#if agent_tunnel.wss_proxy}}
Environment=WSS_PROXY={{agent_tunnel.wss_proxy}}
{{/if}}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
